# image/etc/s6-overlay/s6-rc.d/backend-log/run
#!/command/execlineb -P
logutil-service /data/log/backend
