# image/etc/s6-overlay/s6-rc.d/backend/run
#!/bin/sh
set -a
. /config.cfg
. /data/kafka/env
. /run/secrets/secrets.cfg
exec 2>&1
s6-setuidgid backend /opt/logstash/bin/logstash --path.settings /data/config/backend
