# image/etc/s6-overlay/s6-rc.d/backend/run
#!/command/execlineb -P
if { echo Starting backend }
sleep 1d
