# image/etc/s6-overlay/s6-rc.d/frontend/run
#!/command/execlineb -P
if { echo Starting frontend }
sleep 1d
