# image/etc/s6-overlay/s6-rc.d/kafka/run
#!/command/execlineb -P
if { echo Starting kafka }
sleep 1d
