# image/etc/s6-overlay/s6-rc.d/kafka/run
#!/bin/sh
exec 2>&1
s6-setuidgid kafka /opt/kafka/bin/kafka-run-class.sh kafka.Kafka /data/kafka/server.properties
