# image/etc/s6-overlay/s6-rc.d/log-prepare/up
if { mkdir -p /data/log/backend /data/log/frontend /data/log/kafka }
if { chown nobody:nogroup /data/log/backend /data/log/frontend /data/log/kafka }
chmod 02755 /data/log/backend /data/log/frontend /data/log/kafka
